# debian/rpi-wifi-ap.postinst
#!/bin/bash -e

is_first_install() {
	if [ -z "${2}" ]; then
		return 0
	fi
	return 1
}

if is_first_install "$@"; then
	if [ "$(deb-systemd-helper show -p LoadState --value hostapd)" == "masked" ]; then
        deb-systemd-helper unmask hostapd
    fi
fi

#DEBHELPER#
